# deploy/b70-daily-driver.service
[Unit]
# B70 daily-driver vLLM serve (default: qwen36-27b-int4, DP=2 = two single-card TP=1 replicas, cudagraph=NONE,
# API-key enforced). WEDGE-PROOF: single-card replicas cannot DEVICE_LOST/BCS-wedge, so this survives an
# unattended/traveling weekend (w8a8 TP=2 is faster but wedged under load even post-GuC-fix -- attended-only).
# Managed lifecycle + auto-start on boot. See deploy/README.md and JOURNAL 2026-06-26.
Description=B70 daily-driver vLLM serve
After=docker.service network-online.target mnt-vm_8tb.mount
Requires=docker.service
Wants=network-online.target

[Service]
# oneshot + RemainAfterExit: the daily-driver launches the serve DETACHED (nohup setsid gpu-run + docker run -d)
# and waits for /health, then ExecStart returns 0 and systemd marks the unit active. `systemctl stop` tears it
# down. `systemctl enable` -> starts on boot. (This gives boot-start + clean start/stop; for full crash-restart
# of the container, see the Restart note in deploy/README.md.)
Type=oneshot
RemainAfterExit=yes
User=hotschmoe
Group=hotschmoe
WorkingDirectory=/mnt/vm_8tb/github/b70_ai_things
Environment=DD_MODEL=qwen36-27b-int4
Environment=DD_REPLICAS=2
Environment=DD_MAXLEN=131072
Environment="DD_ENV=GRAPH=1 CGMODE=NONE UTIL=0.95"
# API key is read from the off-repo secret file at start time (NOT baked into the unit).
ExecStart=/bin/bash -lc 'DD_API_KEY="$(cat /mnt/vm_8tb/b70/secrets/dd_api_key)" /mnt/vm_8tb/github/b70_ai_things/daily_driver_serve.sh start'
ExecStop=/bin/bash -lc '/mnt/vm_8tb/github/b70_ai_things/daily_driver_serve.sh stop'
# Two int4 replicas load + JIT-compile (NONE: no capture); allow generous startup.
TimeoutStartSec=1200
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
