# vllm/deploy/b70-dd-watchdog.service
[Unit]
# B70 daily-driver WATCHDOG -- content-aware health/garbage detector + self-heal for the daily-driver serve
# (backend-agnostic: TP=2 single-replica W8A8 default, or the older int4 DP=2 -- set DD_WATCH_REPLICAS to match).
# Probes each replica's OUTPUT (not just /health), and on confirmed "!!!!"/degenerate garbage writes a dated
# incident file (docker-logs grep + dmesg xe/GuC + stats) then `docker restart`s ONLY that replica (no reboot).
# Runs ALONGSIDE the serve (manual launch OR b70-daily-driver.service) -- it does not own the serve lifecycle.
# See bin/dd-watchdog and JOURNAL 2026-06-26 (the ~7h "!!!!" incident this closes the diagnostics gap on).
Description=B70 daily-driver watchdog (content probe + self-heal)
After=docker.service network-online.target mnt-vm_8tb.mount
Wants=network-online.target

[Service]
Type=simple
User=hotschmoe
Group=hotschmoe
WorkingDirectory=/mnt/vm_8tb/github/b70_ai_things
# Knobs (override here if needed): probe every 60s, 2 strikes before heal, 7min cooldown post-restart,
# cap 4 restarts/replica/hr (then ESCALATE-and-log instead of looping a likely-hardware fault).
Environment=DD_WATCH_INTERVAL=60
Environment=DD_WATCH_DEBOUNCE=2
Environment=DD_WATCH_COOLDOWN=420
Environment=DD_WATCH_MAX_RESTARTS_HR=4
Environment=DD_API_KEY_FILE=/mnt/vm_8tb/b70/secrets/dd_api_key
ExecStart=/usr/bin/python3 /mnt/vm_8tb/github/b70_ai_things/bin/dd-watchdog
# Keep the watchdog itself alive (it must outlive transient errors); it does NOT restart on a clean stop.
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
