# vllm/deploy/b70-daily-driver.service
[Unit]
# B70 daily-driver serve. DEFAULT (2026-07-08, user choice): vLLM NVFP4 qwen3.6-27b TP=2 captured+MTP + the
# graph-replay RECLAIM leak fix (B70_XPU_CG_RECLAIM, now the GRAPH=1 default in serve_nvfp4_27b.sh) --
# full-speed ~38 t/s, MTP on, CRASH-FREE (27k single-stream FLAT + 36k 6-way concurrent, 0 linear_stream
# aborts; see JOURNAL 2026-07-08 session 2 + docs/SERVE_TOMORROW_fallback.md). Served as hotschmoe-dd on
# :18080, API-key enforced (key read from the off-repo secret file at start; NOT baked into the unit).
# Fallbacks if the reclaim ever needs disabling: CGRECLAIM=0 (+ B70_EXTRA_ENV=B70_XPU_DRAFTER_EAGER=1 for
# drafter-eager, or drop MTPTOK for B4). dd-watchdog (b70-dd-watchdog.service) backstops SOFT garbage with a
# per-replica docker restart (probes b70_daily_0:18080). Prior default was sglang/vLLM W8A8 27b (see git log).
Description=B70 daily-driver serve (vLLM NVFP4 27b captured+MTP+reclaim, TP=2, served as hotschmoe-dd)
After=docker.service network-online.target mnt-vm_8tb.mount
Requires=docker.service
Wants=network-online.target

[Service]
# oneshot + RemainAfterExit: ExecStart launches the serve container DETACHED (gpu-run -> docker run -d) and
# returns 0; ExecStartPost waits (best-effort) for /health so ordering is sane; `systemctl stop` tears it down.
# `systemctl enable` -> starts on boot. The reclaim makes the container crash-free, so the oneshot pattern is
# sufficient; dd-watchdog handles the rare soft-garbage bounce.
Type=oneshot
RemainAfterExit=yes
User=hotschmoe
Group=hotschmoe
WorkingDirectory=/mnt/vm_8tb/github/b70_ai_things
# Serve config is passed as env to the NVFP4 experiment serve.sh (which reads them). B70_XPU_CG_RECLAIM is NOT
# listed here on purpose -- it is the GRAPH=1 default inside serve_nvfp4_27b.sh (opt out with CGRECLAIM=0).
ExecStart=/bin/bash -lc 'API_KEY="$(cat /mnt/vm_8tb/b70/secrets/dd_api_key)" \
  NAME=b70_daily_0 PORT=18080 TP=2 MODE=fused GRAPH=1 MTPTOK=5 KV_FP8=0 CAPSIZES=1,2,4,8 \
  MAXLEN=131072 MAXSEQS=8 UTIL=0.85 PUSH_AR=1 PUSH_AR_GRAPH=1 SERVED_FORCE=hotschmoe-dd \
  TOOLCALL=1 TOOLPARSER=qwen3_coder REASONPARSER=qwen3 \
  /mnt/vm_8tb/github/b70_ai_things/bin/gpu-run bash /mnt/vm_8tb/github/b70_ai_things/vllm/nvfp4/serve_nvfp4_27b.sh start'
ExecStartPost=/bin/bash -lc 'for i in $(seq 1 120); do curl -sf http://localhost:18080/health >/dev/null 2>&1 && exit 0; docker inspect -f "{{.State.Status}}" b70_daily_0 2>/dev/null | grep -q exited && { echo "container exited during startup" >&2; exit 1; }; sleep 5; done; echo "health wait timed out (container may still be warming)" >&2; true'
ExecStop=/bin/bash -lc 'docker rm -f b70_daily_0 >/dev/null 2>&1 || true'
# NVFP4 load + PIECEWISE capture + MTP drafter: allow generous startup.
TimeoutStartSec=1200
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
